# bin/tool_ashiato_update_log.cron
#! /bin/sh
cd $1
$2 tool_ashiato_update_log.php
