# bin/tool_send_message_in_queue.cron
#! /bin/sh
cd $1
$2 tool_send_message_in_queue.php

