# bin/tool_vacuum_exec.cron
#! /bin/sh
cd $1
$2 -f tool_vacuum_exec.php
