# bin/tool_rss_cache.cron
#! /bin/sh
cd $1
$2 -f tool_rss_cache.php

