# bin/tool_vacuum_exec.cron
#! /bin/sh
cd $1
$2 tool_vacuum_exec.php
