# bin/tool_send_schedule_mail.cron
#! /bin/sh
cd $1
$2 tool_send_schedule_mail.php

