# bin/tool_send_birthday_mail.cron
#! /bin/sh
cd $1
$2 tool_send_birthday_mail.php

