# bin/tool_send_dairy_news.cron
#! /bin/sh
cd $1
$2 tool_send_dairy_news.php

