# bin/tool_ashiato_update_log.cron
#! /bin/sh
cd $1
$2 -f tool_ashiato_update_log.php
