# bin/tool_send_daily_news.cron
#! /bin/sh
cd $1
$2 -f tool_send_daily_news.php
